# testnet-setup/sourcehubd.service
[Unit]
Description=SourceHub service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple

User=ec2-user
Group=ec2-user

ExecStart=/usr/local/bin/sourcehubd start
TimeoutStopSec=5

[Install]
WantedBy=multi-user.target
